# lib/thin/controllers/service.sh.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          thin
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      S 0 1 6
# Short-Description: thin initscript
# Description:       thin
### END INIT INFO

# Original author: Forrest Robertson

# Do NOT "set -e"

### Added Updates
#  Allow Restart of a start,stop and restart of a single config file
#  Author: William Basmayor
#  Example Usage: thin start-config /etc/thin/someconfig.yml 
###

DAEMON=<%= Command.script %>
SCRIPT_NAME=<%= INITD_PATH %>
CONFIG_PATH=<%= config_files_path %>

# Exit if the package is not installed
[ -x "$DAEMON" ] || exit 0

case "$1" in
  start)
	$DAEMON start --all $CONFIG_PATH
	;;
  stop)
	$DAEMON stop --all $CONFIG_PATH
	;;
  restart)
	$DAEMON restart --all $CONFIG_PATH
	;;
  start-config)
  $DAEMON start -C $2
  ;;
  stop-config)
  $DAEMON stop -C $2
  ;;
  restart-config)
  $DAEMON restart -C $2
  ;;
  *)
	echo "Usage: $SCRIPT_NAME {start|stop|restart|start-config|stop-config|restart-config}" >&2
	exit 3
	;;
esac

:
